# templates/ceph-radosgw.service.erb
# cluster related settings are not really overrides if the cluster name is 'ceph'

[Service]
Environment=CLUSTER=<%= @cluster %>
<% if @cpu_shares -%>
CPUShares=<%= @cpu_shares -%>
<% end -%>
<% if @cpu_quota -%>
CPUQuota=<%= @cpu_quota -%>
<% end -%>

